# cachyos-mirrorlist-sync/system/cachyos-mirrorlist-sync.service
[Unit]
Description=Sync CachyOS mirrorlists from upstream PKGBUILDs
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/cachyos-mirrorlist-sync
