# cachyos-mirrorlist-sync/system/cachyos-mirrorlist-sync.timer
[Unit]
Description=Daily CachyOS mirrorlist sync

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
